# Library/LaunchAgents/com.openantigravity.ratelimit-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openantigravity.ratelimit-watcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jeremy/.open-auth-rotator/antigravity/watcher_runner_wrapper.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>EnvironmentVariables</key>
    <dict>
        <key>WORKSPACE_DOMAIN</key>
        <string>zukunftsorientierte-energie.de</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/jeremy</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/antigravity-watcher-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/antigravity-watcher-launchd.log</string>
</dict>
</plist>
